# component/web/src/deb/control/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/visuwall" ]; then
	invoke-rc.d visuwall stop || exit $?
fi
